# .infrastructure/docker-compose/crontab
SHELL=/bin/bash

0 * * * * docker restart proxy-scraper
0 */2 * * * docker exec php-fpm bin/console app:parse:publisher
0 */24 * * * docker exec php-fpm bin/console app:parse:position
